# pi/crontab
@reboot sudo /home/pi/examples/Adafruit_Python_DHT/examples/broadcastTempData.py

